# speedtest.cron
58,18,38 * * * * root /usr/local/speedtest-cli/speedtest-cli --simple > /tmp/speedtest.out
